# scripts/rubymik-demo-reset.timer
[Unit]
Description=Nightly RubyMIK demo reset

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
